# br_external/board/raspberrypi3/rootfs_overlay/etc/crontabs/root
# Rotate logs daily
0 0 * * * /usr/sbin/logrotate
