# br_external/board/raspberrypi3/rootfs_overlay/etc/crontabs/root
# Compress old logs daily
0 0 * * * /sbin/compress-direwolf-logs.sh
