# br_external/board/raspberrypi3/rootfs_overlay/etc/profile.d/less_options.sh
# Enable parsing ANSI color codes
export LESS=-R
